# editions/desktop/live/usr/lib/systemd/system/plasma-log-capture.service
[Unit]
Description=Make ttyS1 writable so the live Plasma session can log to it (diagnostic)
DefaultDependencies=no
After=dev-ttyS1.device
Before=getty@tty1.service
[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/bin/sh -c 'chmod 0666 /dev/ttyS1 2>/dev/null || true'
[Install]
WantedBy=multi-user.target
